# cron_jobs/cron.daily
@daily /home/consularparadi/.dotfiles/scripts/git-sync
